# profile.d/goToFolder.config
alias,folder
s,~\Sources
d,~\Downloads
p,~\Pictures
tpd,"~\Sources\Inexto Suite\Inextrack\main"