# profile.d/goToFolder.config
alias,folder
s,D:\Sources
d,~\Downloads
p,~\Pictures
tpd,"~\Sources\Inexto Suite\Inextrack\main"