# profile.d/ls.alias.ps1
# ls
Set-Alias ls Get-ChildItem-Color -option AllScope -Force
# dir
Set-Alias dir Get-ChildItem-Color -Option AllScope -Force
